# install/cron.d/update-quota
#
# cron.d/update-quota -- synchronize xfs_quota into database
#

*/5 * * * root update-quoga.pl
